# com.minibot.supervisor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.minibot.supervisor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/minibot/projects/my-first-repo/dist/host/supervisor.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/minibot/projects/my-first-repo</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>RUNTIME</key>
    <string>docker</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/minibot/projects/my-first-repo/logs/supervisor.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/minibot/projects/my-first-repo/logs/supervisor.err.log</string>

  <key>ThrottleInterval</key>
  <integer>3</integer>
</dict>
</plist>
